# plain-supervisor.conf
[program:plain]
command = /usr/bin/python3 -m flask run --host=0.0.0.0 --port=8091
environment=FLASK_APP=plain.py
directory=/root/plain
autorestart=true
